# etc/service/k8s-service-ingress-port-bind-reconciler/run
#!/bin/bash
exec /usr/local/bin/k8s-service-ingress-port-bind-reconciler.sh
